# training/deploy/systemd/rfcai-ingestion-daemon.service
[Unit]
Description=RF Connector AI ingestion daemon
After=network-online.target
Wants=network-online.target
# Prevent crash-loop: max 5 restarts in 60s, otherwise leave service in failed.
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
Type=exec
User=rfcai
Group=rfcai
WorkingDirectory=/opt/rfcai/training
EnvironmentFile=/etc/default/rfcai-training
ExecStart=/opt/rfcai/training/.venv/bin/python -m scripts.ingestion_daemon \
    --incoming-dir ${RFCAI_LOCAL_INCOMING} \
    --labeled-root ${RFCAI_LABELED_ROOT} \
    --quarantine-root ${RFCAI_QUARANTINE_ROOT} \
    --classifier-dir ${RFCAI_LOCAL_MODEL_DIR} \
    --interval 5
Restart=on-failure
RestartSec=5

NoNewPrivileges=yes
ProtectSystem=strict
# ProtectHome=yes hides all of /home from the service even with
# ReadWritePaths=, blocking access to the rfcai user's own home dir.
# Keep it off; ProtectSystem + ReadWritePaths is enough scoping.
ReadWritePaths=/home/rfcai
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
