# training/deploy/systemd/rfcai-predict.service
[Unit]
Description=RF Connector AI prediction service (detect+crop+classify)
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
Type=exec
User=rfcai
Group=rfcai
WorkingDirectory=/opt/rfcai/training
EnvironmentFile=/etc/default/rfcai-predict
ExecStart=/opt/rfcai/training/.venv/bin/uvicorn \
    rfconnectorai.server.predict_service:app \
    --host 127.0.0.1 \
    --port 8503 \
    --workers 1 \
    --proxy-headers \
    --forwarded-allow-ips=127.0.0.1
Restart=on-failure
RestartSec=5

NoNewPrivileges=yes
ProtectSystem=strict
ReadWritePaths=/home/rfcai
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
